# etc/service/aesm/run
#!/bin/bash -e

exec chpst -u aesmd:aesmd /opt/intel/sgxpsw/aesm/aesm_service --no-daemon 2>&1
